# ops/systemd/openclaw-ai-tools-digest.service
[Unit]
Description=OpenClaw AI Tools Daily Digest
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/openclaw/openclaw.env
ConditionPathExists=/opt/clawdio/scripts/research_flow_runtime.py

[Service]
Type=oneshot
WorkingDirectory=/opt/clawdio
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /opt/clawdio/scripts/research_flow_runtime.py --env-file /etc/openclaw/openclaw.env run --workflow ai_tools_watch --apply --json

[Install]
WantedBy=default.target
