# ops/systemd/openclaw-job-search-report.service
[Unit]
Description=OpenClaw Job Search Daily Report
After=network-online.target
Wants=network-online.target
ConditionPathExists=/etc/openclaw/openclaw.env
ConditionPathExists=/opt/clawdio/scripts/job_search_assistant.py

[Service]
Type=oneshot
WorkingDirectory=/opt/clawdio
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/env python3 /opt/clawdio/scripts/job_search_assistant.py publish-report --env-file /etc/openclaw/openclaw.env --apply

[Install]
WantedBy=default.target
