# ops/systemd/openclaw-gmail-processor.timer
[Unit]
Description=Run OpenClaw Gmail Inbox Processor every 15 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=15min
AccuracySec=1min
Persistent=true
Unit=openclaw-gmail-processor.service

[Install]
WantedBy=timers.target
